# config/glite-wms-lm.in
#!/bin/bash
#
# glite-wms-lm:      Starts the LogMonitor daemon
#
# Version:        @(#) /etc/rc.d/init.d/glite-wms-lm   2.0
#
# chkconfig: 345 95 06 
# description: Starts, stops and checks the EGEE \
#              LogMonitor daemon.
#
# processname: glite-wms-log_monitor
# config: /opt/glite/etc/glite_wms.conf
# hide: false

# Source function library
. /etc/rc.d/init.d/functions

GLITE_WMS_LOCATION=${GLITE_WMS_LOCATION:-@prefix@}

. ${GLITE_WMS_LOCATION}/etc/glite-wms-vars.sh
. ${GLITE_WMS_LOCATION}/etc/profile.d/glite-wms.sh
. ${GLITE_WMS_LOCATION}/etc/profile.d/glite-wms-config.sh

init_variables()
{
    # All these things must go in some centralized place...
    GLITE_WMS_BIN_DIR=${GLITE_WMS_LOCATION}/bin

    # All these things should stay there
    CONFIGURATION_FILE=glite_wms.conf

    MONITORBASE=glite-wms-log_monitor
    LOGMONITOR=${GLITE_WMS_BIN_DIR}/${MONITORBASE}

    LOCKFILE=`${GLITE_WMS_BIN_DIR}/glite-wms-get-configuration LogMonitor.LockFile`
    SCRIPT_UID=`/usr/bin/id -u`

    GLITEUSER_HOME=`eval echo ~${GLITE_WMS_USER}`
    [ -d ${GLITEUSER_HOME} ] || do_failure "Missing user directory ${GLITEUSER_HOME}"
}

create_dirs()
{
  CONFIGURATION_FILE=glite_wms.conf

  . ${GLITE_WMS_LOCATION}/libexec/glite-wms-parse-configuration.sh ${CONFIGURATION_FILE} LogMonitor
}

remove_lockfile()
{
  if [ -f "${LOCKFILE}" ]; then
    action "A stale lock file still exists, removing it.." /bin/rm -f  ${LOCKFILE}
  fi
}


do_failure()
{
    echo -n " $*"
    failure $1
    echo ""

    exit 1
}

create_proxy()
{
    local proxy=

    if [ ${SCRIPT_UID} -eq 0 ]; then
	proxy=`su $GLITE_WMS_USER -s /bin/sh -c "${GLOBUS_LOCATION}/bin/grid-proxy-init -q \
	    -cert $X509_HOST_CERT \
	    -key  $X509_HOST_KEY \
	    -hours 24 \
	    -out  ${X509_USER_PROXY}.$$ 2>&1"`

	[ $? = 0 ] || do_failure $proxy
    else
	proxy=`${GLOBUS_LOCATION}/bin/grid-proxy-init -q \
	    -cert $X509_HOST_CERT \
	    -key  $X509_HOST_KEY \
	    -hours 24 \
	    -out  ${X509_USER_PROXY}.$$ 2>&1`

	[ $? -eq 0 ] || do_failure $proxy
    fi

    local move=`mv ${X509_USER_PROXY}.$$ ${X509_USER_PROXY} 2>&1`
    [ $? -eq 0 ] || do_failure $move
}

start()
{
    daemon ${LOGMONITOR} -c ${CONFIGURATION_FILE}
    local result=$?
    echo ""

    if [ ${result} -eq 1 ]; then # Startup of the daemon is failed, try to understand the cause
	local pid=`pidofproc ${LOGMONITOR}`

	if [ -z "${pid}" -a -f "${LOCKFILE}" ]; then
	    echo "LogMonitor is not running, but a stale lock file exists."
	    echo "Check situation and try to start again."
	    echo "Lock file path is: ${LOCKFILE}"
	elif [ -n "${pid}" ]; then
	    echo "LogMonitor yet running in pid ${pid}"
	fi
    elif [ ${SCRIPT_UID} -eq 0 ]; then
	echo `pidofproc ${LOGMONITOR}` > /var/run/${MONITORBASE}.pid
    fi
}

reload()
{
    local pid=`pidofproc ${LOGMONITOR}`

    if [ -z "${pid}" -a -f "${LOCKFILE}" ]; then
	pid=`/bin/cat ${LOCKFILE}`
    fi

    if [ -n "${pid}" ]; then
	action "Reloading LogMonitor configuration " /bin/kill -HUP ${pid}
    fi
}

stop()
{
    local times pidfile pid lpid result=0

    pid=`pidofproc ${LOGMONITOR}`
    if [ -f "/var/run/${MONITORBASE}.pid" ]; then
        pidfile=`/bin/cat /var/run/${MONITORBASE}.pid`
    fi

    if [ -n "${pid}" -a -n "${pidfile}" ]; then
        if [ "${pid}" != "${pidfile}" ]; then
            result=1
            echo -ne "\t\tCould not reliably find LogMonitor pid!\n"
            pid=
            pidfile=
        fi
    else
        if [ -z "${pidfile}" ]; then
            pidfile=$pid
        fi
    fi

    if [ -n "${pidfile}" ]; then
        if [ -f "${LOCKFILE}" ]; then
           lpid=`/bin/cat ${LOCKFILE}`
        fi
        if [ "${pidfile}" == "${lpid}" ]; then
            if [ -n "${pid}" ]; then
                /bin/kill -TERM $pid
                for (( times = 10; times >= 0; times-- )); do
                    [ -f "${LOCKFILE}" ] || { success $"LogMonitor terminated normally" && echo && break; }
                    if [ $times -eq 0 ]; then
                        killproc ${LOGMONITOR}
                        result=$?
                        echo
                    else
                        sleep 1
                    fi
                done
                if [ $result -eq 0 ]; then
                    remove_lockfile
                    if [ ${SCRIPT_UID} -eq 0 ]; then
                        rm -f /var/run/${MONITORBASE}.pid
                    fi
                fi
            else
                result=1
                failure $"LogMonitor not running, but lock file found"
                echo
            fi
        else
            if [ ! -f "${LOCKFILE}" ]; then
                if [ -n "${pid}" ]; then
                    action "" /bin/kill -KILL $pid
                    result=$?
                else
                    success $"LogMonitor was not running"
                    echo
                fi
                if [ $result -eq 0 -a ${SCRIPT_UID} -eq 0 ]; then
                    rm -f /var/run/${MONITORBASE}.pid
                fi
            else
                result=1
                failure $"Inconsistent LogMonitor lock files"
                echo
            fi
        fi
    else
        if [ $result -eq 0 ]; then
            if [ -f "${LOCKFILE}" ]; then
                success
                echo
                remove_lockfile
            else
                echo -ne "\t\tLogMonitor not running!\n"
            fi
        fi
    fi
}

status()
{
    local pid=

    pid=`pidofproc ${LOGMONITOR}`

    if [ -n "$pid" ]; then
	echo "Logmonitor running..."
    elif [ -f "${LOCKFILE}" ]; then
	echo "LogMonitor not running but stale lock file present."
	exit 1
    else
	echo "LogMonitor stopped."
    fi
}

check()
{
    status=`${LOGMONITOR} -Cc ${CONFIGURATION_FILE} 2>&1`

    if [ $? -eq 0 ]; then
	success $"check"
	echo ""
    else
	failure $"check"
	echo ""
	echo $status
    fi
}

cd /tmp

case $1 in
    start)
	echo -n "Starting LogMonitor..."

	init_variables
        create_dirs	
	start
	RETVAL=$?
    ;;
    stop)
	echo -n "Stopping LogMonitor..."

	init_variables
	stop
	RETVAL=$?
    ;;
    restart)
	echo -n "Stopping LogMonitor..."

	init_variables
	create_dirs
	stop

	echo -n "Starting LogMonitor..."
	start
	RETVAL=$?
    ;;
    reload)
	init_variables
	reload
	RETVAL=$?
    ;;
    status)
	init_variables
	status
	RETVAL=$?
    ;;
    check)
        echo -n "Checking installation..."

	init_variables
	check
	RETVAL=$?
    ;;
    *)
	echo "Usage: $0 {start|stop|restart|status|check}"
	RETVAL=1
    ;;
esac
exit $RETVAL
