# src/edg-wl-purgeStorage-hourly.crontab
# Execute the "purger" command at every day except on Sunday with a frequency of one hour
# if and only if the percentage of allocated blocks is greater than 40%
0 */1 * * mon-sat $EDG_WL_LOCATION/sbin/edg-wl-purgeStorage -l $EDG_WL_LOCATION_VAR/log/edg-wl-purgeStorage.log -t 604800 -a 40
